# simulation-service-prototype/etc/systemd/system/simulator-cleanup.service
[Unit]
Description=Clean up PX4 and dronology logs

[Service]
Type=simple
ExecStart=/usr/local/bin/simulator-cleanup